# telegram-bot-api.service.example
[Unit]
Description=Telegram Bot API Local Server (Docker)
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10
EnvironmentFile=/etc/telegram-bot/project.env
EnvironmentFile=/etc/telegram-bot/api.env
ExecStartPre=-/usr/bin/docker rm -f telegram-bot-api
ExecStart=/usr/bin/docker run --rm --name telegram-bot-api \
    --network host \
    -e TELEGRAM_API_ID=${TELEGRAM_API_ID} \
    -e TELEGRAM_API_HASH=${TELEGRAM_API_HASH} \
    -e TELEGRAM_LOCAL=true \
    -e TELEGRAM_STAT=true \
    -v ${PROJECT_DIR}/telegram-bot-api-data:/var/lib/telegram-bot-api \
    -v ${PROJECT_DIR}/video_downloads:${PROJECT_DIR}/video_downloads:ro \
    aiogram/telegram-bot-api:latest
ExecStop=/usr/bin/docker stop telegram-bot-api

[Install]
WantedBy=multi-user.target
