# telegram_bot.service
[Unit]
Description=Telegram Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=huzhuoliang
WorkingDirectory=/home/huzhuoliang/claude-workspace/telegram_monitor
ExecStart=/usr/bin/python3 /home/huzhuoliang/claude-workspace/telegram_monitor/bot.py
Environment=PATH=/home/huzhuoliang/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# To use claude_backend="api", create /etc/telegram_bot.env (chmod 600) with:
#   ANTHROPIC_API_KEY=sk-ant-...
# Then uncomment the line below:
# EnvironmentFile=/etc/telegram_bot.env

[Install]
WantedBy=multi-user.target
